# supervisor.conf
[supervisord]
nodaemon=false
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 1440x900x24 +extension GLX +render
autostart=true
autorestart=true
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.log

[program:fluxbox]
command=fluxbox
autostart=true
autorestart=true
stdout_logfile=/data/logs/fluxbox.log
stderr_logfile=/data/logs/fluxbox.log

[program:x11vnc]
command=x11vnc -display :99 -forever -shared -rfbport 5900 -nopw
autostart=true
autorestart=true
stdout_logfile=/data/logs/vnc.log
stderr_logfile=/data/logs/vnc.log

[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
stdout_logfile=/data/logs/novnc.log
stderr_logfile=/data/logs/novnc.log

[program:api]
command=python3 /app/api/server.py
autostart=true
autorestart=true
stdout_logfile=/data/logs/api.log
stderr_logfile=/data/logs/api.log
